# ansible/templates/cfwarp-dogfood-observer.timer.j2
[Unit]
Description=Run cfwarp dogfood observer probe periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec={{ dogfood_observer_interval }}
AccuracySec=30s
Unit=cfwarp-dogfood-observer.service
Persistent=true

[Install]
WantedBy=timers.target
